# config/birdcam-tcp.service
[Unit]
Description=Piepswatch Bird Camera TCP Stream
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/emily/piepswatch-birdcam/scripts/birdcam-tcp-stream.sh
Restart=always
RestartSec=5
User=emily
Group=emily
StandardOutput=journal
StandardError=journal
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
